# deploy/backtest-backend.service
[Unit]
Description=BackTest Studio API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/credit-backtest-studio/backend
EnvironmentFile=-/var/www/credit-backtest-studio/backend/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
